# dotconfig/systemd/user/timers.target.wants/trash-empty.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily
Persistent=true

[Unit]
Description=Run trash-empty daily
